# etc/sv/sshd/run
#!/bin/sh

exec /usr/sbin/sshd -De 2>&1
