# etc/sv/ngetty/run
#!/bin/sh

exec chpst -e ./env /sbin/ngetty tty2 tty3 tty4
