# etc/sv/kdm/run
#!/bin/sh

exec /usr/sbin/kdm -nodaemon 1>&2
