# etc/sv/agetty-generic/run
#!/bin/sh
service=$(basename $(pwd))
tty=${service/*-}

exec /sbin/agetty -8 38400 $tty linux
