# etc/sv/mongodb/run
#!/bin/sh

exec chpst -u mongodb:daemon /usr/bin/mongod -f /etc/mongodb.conf 2>&1
