# etc/sv/openntpd/run
#!/bin/sh

exec /usr/sbin/ntpd -s -d 2>&1
