# etc/sv/apache2/run
#!/bin/sh
exec 2>&1
exec httpd -eDEBUG -DNO_DETACH -DFOREGROUND
