# etc/sv/syslog-ng/finish
#!/bin/sh
exec 2>&1
sed -e -i 's|unix-dgram("/dev/log");|unix-dgram("/run/systemd/journal/syslog");|' /etc/syslog-ng/syslog-ng.conf
