# etc/sv/syslog-ng/run
#!/bin/sh
exec 2>&1
sed -i -e 's|unix-dgram("/run/systemd/journal/syslog");|unix-dgram("/dev/log");|' /etc/syslog-ng/syslog-ng.conf
exec syslog-ng -F
