# etc/sv/bumblebeed/run
#!/bin/sh

exec /usr/sbin/bumblebeed 2>&1
