# etc/sv/couchdb/run
#!/bin/sh

export HOME=/var/lib/couchdb

exec chpst -u couchdb:daemon /usr/bin/couchdb 2>&1
