# etc/sv/gdm/run
#!/bin/sh

exec /usr/bin/gdm -nodaemon 1>&2
