# etc/sv/gdm/run
#!/bin/sh

exec /usr/sbin/gdm -nodaemon 1>&2
