# etc/sv/postgresql/run
#!/bin/sh

. /etc/conf.d/postgresql

exec chpst -u postgres /usr/sbin/postgres -D $PGROOT/data -c config_file=$PGROOT/data/postgresql.conf
