# etc/sv/avahi/run
#!/bin/sh
sv start dbus || exit 1
exec avahi-daemon
